# ansible/roles/thunderbolt_net/templates/10-tb-up.rules.j2
# Managed by Ansible (thunderbolt_net).
ACTION=="add", SUBSYSTEM=="net", NAME=="{{ storage_iface }}", RUN+="/usr/local/sbin/tb-link-up.sh {{ storage_iface }}"
