# ansible/roles/thunderbolt_net/templates/tb-link.timer.j2
[Unit]
Description=ailab Thunderbolt storage link reconciler

[Timer]
OnBootSec=10
OnUnitActiveSec=30
AccuracySec=5

[Install]
WantedBy=timers.target
