# ansible/roles/storage_router/templates/storage-router.service.j2
[Unit]
Description=Storage-router SNAT reconcile (Talos VM -> QNAP fabric over Thunderbolt)
After=network.target
Wants=storage-router.timer

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/storage-router.sh
